# dot_config/profile.d/20-ghcup.sh
if [ -r "$HOME/.ghcup/env" ]; then
    . "$HOME/.ghcup/env"
fi
