# dot_config/profile.d/20-mise.sh
if command -v mise >/dev/null; then
    eval "$(mise env --shell bash)"
fi
