# private_dot_config/profile.d/20-juliaup.sh
case ":$PATH:" in
    *:"$HOME/.juliaup/bin":*)
        # Do nothing
        ;;
    *)
        export PATH="$HOME/.juliaup/bin:$PATH"
        ;;
esac
